# oculink-gpu-monitor.service
[Unit]
Description=OCuLink GPU Hot-plug Monitor
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/oculink-gpu-watcher
Restart=always
RestartSec=5
User=root
Environment="DISPLAY=:0"
Environment="CHECK_INTERVAL=5"
RuntimeDirectory=oculink-hotplug
CacheDirectory=oculink-hotplug
ConfigurationDirectory=oculink-hotplug

[Install]
WantedBy=multi-user.target
